# macosx/emacs/gnu.emacs.daemon.LaunchAtLogin.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>gnu.emacs.daemon.LaunchAtLogin.agent</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Applications/Emacs.app/Contents/MacOS/Emacs</string>
    <string>--daemon</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>ServiceDescription</key>
  <string>Run GNU Emacs Daemon using MacOSX Launchd</string>
  <!--NOTE: pls change to /usr/local/var/log/emacs.log on MacOSX Catalina-->
  <key>StandardOutPath</key>
  <string>/opt/homebrew/var/log/emacs.log</string>
  <key>StandardErrorPath</key>
  <string>/opt/homebrew/var/log/emacs.log</string>
</dict>
</plist>
